# examples/ergo-signer.service
[Unit]
Description=Ergo Transaction Signer (BlockHost)
After=network.target
Before=blockhost-monitor.service

[Service]
Type=simple
User=blockhost
Group=blockhost
ExecStart=/usr/share/blockhost/ergo-signer
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
